# scripts/nvr.service
[Unit]
Description=NVR System Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/nvr-design/server
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
